# templates/default/secuity-monkey-supervisor.ini.erb
[unix_http_server]
file=/tmp/supervisor.sock;

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock;

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisord]
logfile=/tmp/securitymonkey.log
logfile_maxbytes=50MB
logfile_backups=2
loglevel=trace
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200
user=<%= @user%>

[program:securitymonkey]
command=python /home/<%= @user%>/security_monkey/manage.py run_api_server
environment=SECURITY_MONKEY_SETTINGS="<%= @release_path%>/security_monkey/env-config/config-deploy.py"

[program:securitymonkeyscheduler]
command=python <%= @release_path%>/manage.py start_scheduler

directory=<%= @release_path%>
environment=PYTHONPATH='<%= @release_path%>',SECURITY_MONKEY_SETTINGS="<%= @release_path%>/env-config/config-deploy.py"
user=<%= @user%>
autostart=true
autorestart=true
